# config/supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:comfyui]
directory=/notebooks/ComfyUI
command=/bin/bash -lc 'export HF_HOME="${HF_HOME:-/storage/sd-suite/hf_cache}"; export PYTORCH_CUDA_ALLOC_CONF="${PYTORCH_CUDA_ALLOC_CONF:-expandable_segments:True}"; exec python main.py --listen ${COMFYUI_LISTEN_HOST:-0.0.0.0} --port ${COMFYUI_PORT:-8189} ${COMFYUI_EXTRA_ARGS:---disable-auto-launch}'
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/comfyui.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/tmp/comfyui.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
environment=HOME="/home/%(ENV_MAMBA_USER)s",USER="%(ENV_MAMBA_USER)s"

[program:jupyterlab]
command=/bin/bash -lc 'exec jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --NotebookApp.token="" --NotebookApp.password="" --notebook-dir=/notebooks'
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/jupyterlab.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/tmp/jupyterlab.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
environment=HOME="/home/%(ENV_MAMBA_USER)s",USER="%(ENV_MAMBA_USER)s"
