# config/supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:comfyui]
command=/bin/bash -lc 'cd /opt/app/ComfyUI && exec python main.py --listen 0.0.0.0 --port 8189'
directory=/opt/app/ComfyUI
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=20
stopsignal=TERM
stdout_logfile=/tmp/comfyui.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/tmp/comfyui.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3
redirect_stderr=false
environment=PYTHONUNBUFFERED="1",HF_HOME="/storage/sd-suite/hf_cache",HOME="/home/mambauser",PATH="/opt/conda/envs/pyenv/bin:/opt/conda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
user=mambauser
priority=10
